# packaging/clipdashd.service
[Unit]
Description=Clipdash Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/clipdash-daemon
Restart=on-failure
RestartSec=2
Environment=RUST_LOG=info

[Install]
WantedBy=default.target

